# cfg/obfsproxy.service
[Unit]
Description=OBFSProxy

[Service]
ExecStart=/usr/bin/obfsproxy-wrapper

[Install]
WantedBy=multi-user.target
